# com.djailla.subsonic/source/rc.local
#!/bin/sh -e
#
# rc.local
#

case $1 in
    start)
        echo 'SUBSONIC_ARGS="--port=7575"\nSUBSONIC_USER=root\n' > /etc/default/subsonic
        service subsonic start
    ;;
    stop)
        service subsonic stop
    ;;
esac
