# com.djailla.subsonic/source/rc.local
#!/bin/sh -e
#
# rc.local
#

case $1 in
    start)
        printf 'SUBSONIC_ARGS="--port='$RAINBOW_PORT_WEB_UI'"\nSUBSONIC_USER=root\n' > /etc/default/subsonic
        service subsonic start
    ;;
    stop)
        service subsonic stop
    ;;
esac
